# overlay/fs/etc/systemd/system/getty@tty1.service.d/buffyboard.conf
[Service]
ExecStartPost=-/usr/bin/systemctl start buffyboard.service
StandardError=journal
