# overlay/fs/etc/systemd/system/getty@tty1.service.d/buffyboard.conf
[Unit]
Wants=buffyboard.service

[Service]
ExecStartPost=-/usr/bin/systemctl start buffyboard.service
StandardError=journal
