# overlay/fs/etc/update-motd.d/05-mpwrd-header
#!/bin/bash

# Fancy mPWRD Logo

set -o pipefail

echo -n -e "\e[32m"
cat << 'EOF'
                     .                                          
                   :@@+:.             =*#+-                     
                  +@@@@@.           .*@@@@@#.                   
                :*@@@@*.           :@@@@@@@@#=                  
               =@@@@@=           .*@@@@@+@@@@@=                 
             .*@@@@#:           :#@@@@#: .*@@@@#:               
            :@@@@@*            =@@@@@=    .=@@@@#=              
           +@@@@@:           .*@@@@#:       .#@@@@#.            
         :*@@@@*.           -@@@@@*          .=@@@@#=           
        =@@@@@=           .*@@@@@=             =@@@@@*.         
      .*@@@@#:           :#@@@@#:               .#@@@@#=        
     :@@@@@*            =@@@@@=                  .=@@@@@=       
    +@@@@@:           .*@@@@#:                     .#@@@@#.     
  :*@@@@*.           -@@@@@*                        .=@@@@#=    
 =@@@@@=           .*@@@@@=                           =@@@@@#.  
.*#@@#:            :*@@@#:                             .#@@@#-  
   :+                .-=                                .==.    
EOF
echo -e "\e[0m"
cat << 'EOF'
++*****+=:.  ****            :***: :+******=:.    =+**#**++:.   
@@@#**##@@*. +@@@:   :###.   +@@@  -@@@#*##@@@=   *@@@***#@@@*. 
@@@=   .#@@+ .@@@+   #@@@+   @@@-  -@@@.   +@@@.  *@@*    :+@@@:
@@@=   -@@@=  *@@#  =@@#@#: -@@@   -@@@. .-*@@#   *@@*      #@@+
@@@@#@@@@@+   .@@@-.@@+:@@* *@@+   -@@@@@@@@@+:   *@@*      *@@#
@@@*++==-.     +@@#*@@  +@@+@@@    -@@@=-#@@#.    *@@*     :@@@=
@@@=           .@@@@@:  .#@@@@-    -@@@.  +@@@-   *@@*::--*@@@+ 
@@@=            -@@@*    -@@@*     -@@@.   =@@@=  *@@@@@@@@@+:  
:::.            .:::      :::.     .---     :---  .:::--:..     
EOF
